# infra/systemd/matrummet-image.service
[Unit]
Description=Matrummet Image Service
After=network.target

[Service]
Type=simple
User=deploy_recept
Group=deploy_recept
WorkingDirectory=/opt/matrummet/apps/image-service

# Environment file stored outside deployment directory (safe from rsync --delete)
EnvironmentFile=/opt/matrummet/.image-service.env
Environment="NODE_ENV=production"

ExecStart=/usr/local/bin/bun src/index.ts
Restart=on-failure
RestartSec=5s
StartLimitInterval=300s
StartLimitBurst=10

TimeoutStopSec=15s
KillMode=mixed
KillSignal=SIGTERM

StandardOutput=journal
StandardError=journal
SyslogIdentifier=matrummet-image

[Install]
WantedBy=multi-user.target
